# systemd/tool-snapshot.timer
[Unit]
Description=Toolmaster daily snapshot timer

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
